# packaging/arasul-platform/etc/systemd/system/arasul-platform.service
[Unit]
Description=Arasul Platform - Edge AI Appliance
Documentation=https://github.com/arasul/platform
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=arasul
Group=arasul
WorkingDirectory=/opt/arasul

# Environment - load from config or .env
EnvironmentFile=-/opt/arasul/.env

# Start command - ordered startup for post-reboot reliability
ExecStart=/opt/arasul/scripts/system/ordered-startup.sh --skip-pull
ExecStop=/usr/bin/docker compose down

# Restart policy
Restart=on-failure
RestartSec=30s
TimeoutStartSec=1500
TimeoutStopSec=120

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security
PrivateTmp=yes
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/arasul

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arasul-platform

[Install]
WantedBy=multi-user.target
