# packaging/arasul-platform/etc/systemd/system/deadman-switch.service
[Unit]
Description=ARASUL Deadman Switch for Self-Healing Agent
After=docker.service arasul-platform.service
Wants=arasul-platform.service

[Service]
Type=oneshot
ExecStart=/arasul/scripts/system/deadman-switch.sh
TimeoutStartSec=30
